# privileges-demote/payload/Library/LaunchAgents/corp.sap.privileges.demote.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>corp.sap.privileges.demote</string>
	<key>ProgramArguments</key>
	<array>
    	<string>/bin/zsh</string>
    	<string>-c</string>
    	<string>! log show --last 5m --style syslog --predicate &apos;process == &quot;corp.sap.privileges.helper&quot; &amp;&amp; eventMessage CONTAINS &quot;has now admin&quot;&apos; | grep &quot;corp.sap.privileges.helper&quot; &amp;&amp; echo &quot;user is an admin and did not request privileges within the last 10 minutes, revoking right&quot; &amp;&amp; /Applications/Privileges.app/Contents/Resources/PrivilegesCLI --remove</string>
  	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>600</integer>
</dict>
</plist>
